# bin/dask-worker.service
[Unit]
Description=Private server
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/
Environment="PATH=/home/ubuntu/pp_client/venv/bin:/usr/bin"
ExecStart=/bin/bash -a -c 'source /home/ubuntu/pp_client/venv/bin/activate && dask-worker localhost:8786 --no-nanny --nthreads 2 --resources "process=1" --worker-port WORKER_PORT'
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target

